# reporter-engine.service
[Unit]
Description=Reporter Engine - Test Reporting Platform
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/reporter-engine
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reporter-engine

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/azureuser/reporter-engine

[Install]
WantedBy=multi-user.target
